# traktoros.service
[Unit]
Description=Traktoros Flask API
After=network.target

[Service]
User=ec2-user
WorkingDirectory=/home/ec2-user/TraktorOS
EnvironmentFile=/home/ec2-user/.traktoros-env
ExecStart=/usr/bin/python3 local_server.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
